# supervisord.conf
[inet_http_server]
port=%(ENV_INET_HTTP_SERVER_HOST)s:%(ENV_INET_HTTP_SERVER_PORT)s
chmod=0777
chown=nobody:nogroup
username=appuser
password=wasteantadmin@2024

[supervisord]
nodaemon=true
logfile=/var/log/supervisor/supervisord.log
pidfile=/var/run/supervisord.pid
user=root

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface
  
[supervisorctl]
serverurl=http://%(ENV_INET_HTTP_SERVER_HOST)s:%(ENV_INET_HTTP_SERVER_PORT)s
username=appuser
password=wasteantadmin@2024

##################################################
########### Data Acquisition Module ##############
[program:data_acquisition_endpoint]
environment=PYTHONPATH=/home/%(ENV_user)s/src
command=/bin/bash -c "source /opt/ros/%(ENV_ROS_DISTRO)s/setup.bash && source ~/devel/setup.bash && python3 data_reader/main.py"
directory=/home/%(ENV_user)s/src
autostart=false
autorestart=true
stderr_logfile=/var/log/data_acquisition_endpoint.err.log
stdout_logfile=/var/log/data_acquisition_endpoint.out.log

[program:data_acquisition_server]
environment=PYTHONPATH=/home/%(ENV_user)s/src
command=/bin/bash -c "source /opt/ros/%(ENV_ROS_DISTRO)s/setup.bash && source ~/devel/setup.bash && python3 data_reader/interface/grpc/grpc_server.py"
directory=/home/%(ENV_user)s/src
autostart=true
autorestart=true
stderr_logfile=/var/log/data_acquisition_server.err.log
stdout_logfile=/var/log/data_acquisition_server.out.log


####################################################
############### Segmentation Module ################
[program:segments_data_reader]
environment=PYTHONPATH=/home/%(ENV_user)s/src
command=/bin/bash -c "python3 segments/data_reader/main.py"
directory=/home/%(ENV_user)s/src
autostart=true
autorestart=true
stderr_logfile=/var/log/segments_data_reader.err.log
stdout_logfile=/var/log/segments_data_reader.out.log

[program:segments_comuting_server]
environment=PYTHONPATH=/home/%(ENV_user)s/src
command=/bin/bash -c "python3 segments/interface/grpc/grpc_server.py"
directory=/home/%(ENV_user)s/src
autostart=true
autorestart=true
stderr_logfile=/var/log/segments_comuting_server.err.log
stdout_logfile=/var/log/segments_comuting_server.out.log

######################################################
############### Visualization Module #################
[program:visualizer_data_reader]
environment=PYTHONPATH=/home/%(ENV_user)s/src
command=/bin/bash -c "python3 visualizes/data_reader/main.py"
directory=/home/%(ENV_user)s/src
autostart=true
autorestart=true
stderr_logfile=/var/log/visualizer_data_reader.err.log
stdout_logfile=/var/log/visualizer_data_reader.out.log

[program:visualizer_comuting_server]
environment=PYTHONPATH=/home/%(ENV_user)s/src
command=/bin/bash -c "source /opt/ros/%(ENV_ROS_DISTRO)s/setup.bash && source ~/devel/setup.bash && python3 interface/grpc/grpc_server.py"
directory=/home/%(ENV_user)s/src/visualizes
autostart=true
autorestart=true
stderr_logfile=/var/log/visualizer_comuting_server.err.log
stdout_logfile=/var/log/visualizer_comuting_server.out.log


######################################################
############### Impurity Module #################
[program:impurity_data_reader]
environment=PYTHONPATH=/home/%(ENV_user)s/src
command=/bin/bash -c "python3 impurity/data_reader/main.py"
directory=/home/%(ENV_user)s/src
autostart=true
autorestart=true
stderr_logfile=/var/log/impurity_data_reader.err.log
stdout_logfile=/var/log/impurity_data_reader.out.log

[program:impurity_comuting_server]
environment=PYTHONPATH=/home/%(ENV_user)s/src
command=/bin/bash -c "python3 impurity/interface/grpc/grpc_server.py"
directory=/home/%(ENV_user)s/src
autostart=true
autorestart=true
stderr_logfile=/var/log/impurity_comuting_server.err.log
stdout_logfile=/var/log/impurity_comuting_server.out.log


###############################################
############ Grouping processor ###############
[group:data_acquisition]
programs=data_acquisition_server,data_acquisition_endpoint

[group:segmentation]
programs=segments_data_reader,segments_comuting_server

[group:impurity]
programs=impurity_data_reader,impurity_comuting_server

[group:visualizer]
programs=visualizer_data_reader,visualizer_comuting_server
